# poweroff120.conf
description	"Synology poweroff"

start on runlevel 0 and stopped umount-root-fs and umount-root-ok
stop on runlevel [!0]

task

console none

script
	## make sure runlevel is not 6 (reboot)
	run_level=`runlevel | awk '{ printf $2 }'` || true
	if [ "x${run_level}" = "x6" ]; then
		echo "incorrect runlevel, skip poweroff"
		exit 0
	fi

	if [ "$INIT_HALT" = "" ]; then
		INIT_HALT=POWEROFF
	fi

	# If INIT_HALT=HALT don't poweroff.
	poweroff="-p"
	if [ "$INIT_HALT" = "HALT" ]; then
		poweroff=""
	fi

	i2cset -y -f 0 0x45 0x00 0x55
	i2cset -y -f 0 0x45 0x01 0x01
	i2cset -y -f 0 0x45 0x31 0x33
	i2cset -y -f 0 0x45 0x32 0x33
	i2cset -y -f 0 0x45 0x33 0x33
	i2cset -y -f 0 0x45 0x30 0x07
	i2cset -y -f 0 0x45 0x34 128
	i2cset -y -f 0 0x45 0x35 0
	i2cset -y -f 0 0x45 0x36 0
	
	halt -f $poweroff

end script
# vim:ft=upstart
